# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root nginx -g "daemon off;"