# rootfs/etc/services.d/nas-tools/run
#!/usr/bin/with-contenv bash

cd ${WORKDIR}

umask "${UMASK}"

exec s6-setuidgid "${PUID}":"${PGID}" "$(which dumb-init)" "$(which pm2-runtime)" start run.py -n NAStool --interpreter python3